# etc/systemd/system/vibed.service
# vibed — JSON-RPC daemon that wraps the `vibe` CLI for the admin web app.
#
# Installed by install.sh as /etc/systemd/system/vibed.service. To upgrade
# the binary, install the new /usr/local/bin/vibed and run:
#   sudo systemctl restart vibed

[Unit]
Description=Vibe appliance daemon (JSON-RPC over /run/vibed.sock)
After=network-online.target docker.service
Wants=network-online.target
# Don't block boot if vibed has trouble; the CLI keeps working.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/vibed --socket /run/vibed.sock --vibe-bin /usr/local/bin/vibe
Restart=on-failure
RestartSec=2

# Run as root. The `vibe` CLI vibed shells out to runs require_root for
# installs/upgrades/uninstalls (it chowns per-app data dirs to specific
# uids — postgres 70, redis 999, vibe-user — which non-root can't do).
# We considered running vibed as `vibe` and sudoing per-call, but that
# would force NoNewPrivileges off and add a sudoers config, which is a
# wash on the security side and adds operator surprise. Hardening lives
# in the directives below — no setuid binaries, no kernel-tunable
# writes, the daemon only writes inside the /etc/vibe + /var/lib/vibe
# + /var/log/vibe + /run trees.

# Hardening — even running as root, the namespace restrictions below
# stop a daemon compromise from rewriting random parts of the host.
ProtectSystem=strict
ProtectHome=true
# /etc/vibe is writable because lib/apps.sh's secrets_set / config_set
# write env files there. /opt/vibe-installer stays read-only — the
# installer is the only thing that updates the repo.
ReadWritePaths=/run /var/lib/vibe /var/log/vibe /etc/vibe
ReadOnlyPaths=/opt/vibe-installer
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
# vibed itself doesn't need outbound network — the RPC client connects
# via Unix socket. But the wrapped `vibe` CLI does talk to docker /
# GHCR / licensing.kisaes.com etc., so we leave network access on.

[Install]
WantedBy=multi-user.target
